# 26.04template/crontab/tony_everyday.sh
#!/bin/bash

## Executes at 1 minute past 1:00 AM: 1 1 * * *