# 26.04template/crontab/tony_every_hour.sh
#!/bin/bash

## Executes on the hour every hour: 0 * * * *

## Clear darija attempts
    # /usr/bin/php /var/local/crontab/clear_darija_attempts.php